# com.user.raid-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  LaunchAgent for com.user.raid-monitor
  Template: INSTALL_PATH and USERNAME are replaced by install.sh
  Install location: ~/Library/LaunchAgents/com.user.raid-monitor.plist

  After editing this file manually, reload with:
    launchctl unload  ~/Library/LaunchAgents/com.user.raid-monitor.plist
    launchctl load    ~/Library/LaunchAgents/com.user.raid-monitor.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.raid-monitor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>INSTALL_PATH/raid-monitor.sh</string>
    </array>

    <!-- Fire immediately on load, then every 300 seconds (5 minutes).
         If you change this, also update POLL_INTERVAL_SECONDS in config.sh. -->
    <key>StartInterval</key>
    <integer>300</integer>

    <key>RunAtLoad</key>
    <true/>

    <!-- Explicit environment: launchd does not guarantee HOME or PATH are set.
         ~/bin and /usr/sbin (for diskutil) must both be on PATH.
         Both Homebrew prefixes are included for msmtp / other optional tools. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/USERNAME</string>
        <key>PATH</key>
        <string>/Users/USERNAME/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <!-- Capture stdout and stderr to separate files for easy debugging.
         Logs in the same data directory as the main log. -->
    <key>StandardOutPath</key>
    <string>/Users/USERNAME/.local/share/raid-monitor/raid-monitor-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/USERNAME/.local/share/raid-monitor/raid-monitor-stderr.log</string>
</dict>
</plist>
